# ansible/roles/project/files/websockets.conf
[program:websockets]
command=/usr/bin/php /var/www/vhosts/postbin/live/artisan websockets:serve
numprocs=1
autostart=true
autorestart=true
#user=laravel-echo
user=root